# data/beadm-clear-boot-once.service.in
[Unit]
Description=Clear temporary boot environment activations
Requires=zfs.target
After=zfs.target
ConditionPathExists=/dev/zfs

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@bindir@/beadm --client=libzfs activate -T

[Install]
WantedBy=multi-user.target
