# data/beadm.service.in
[Unit]
Description=Boot Environment Administration Service
ConditionPathExists=/dev/zfs

[Service]
BusName=@bus_name@
ExecStart=@bindir@/beadm daemon
ExecReload=busctl call @bus_name@ /ca/kamacite/BootEnvironments ca.kamacite.BootEnvironmentManager Refresh
Restart=on-failure
Environment=RUST_LOG=info
